# root/etc/cron.weekly/bleachbit
#!/bin/bash 


bleachbit --clean --all-but-warning

sudo bleachbit --clearn --all-but-warning
