# root/etc/cron.hourly/reflector
#!/bin/bash

reflector --latest 20 --protocol http --sort age --save /etc/pacman.d/mirrorlist